# watcher/crosstalk-notify.path
# crosstalk-notify.path
# Watches inbox.jsonl for new messages and triggers the email notification service.
# Install:  cp this to ~/.config/systemd/user/
# Enable:   systemctl --user enable --now crosstalk-notify.path

[Unit]
Description=Watch crosstalk inbox for new messages

[Path]
PathChanged=%h/.crosstalk/inbox.jsonl

[Install]
WantedBy=default.target
