# watcher/crosstalk-notify.service
# crosstalk-notify.service
# Sends email when new crosstalk messages arrive.
# Triggered by crosstalk-notify.path.
# Install:  cp this to ~/.config/systemd/user/
# Enable:   systemctl --user enable --now crosstalk-notify.path

[Unit]
Description=Send email notification for new crosstalk messages
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.crosstalk/watcher/notify.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crosstalk-notify

# Security hardening (runs as you, with minimal access)
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.crosstalk/watcher
ReadOnlyPaths=%h/.crosstalk/inbox.jsonl
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtimeIO=yes
MemoryDenyWriteExecute=yes
